# code/despliegue-plataforma/local/prometheus/roles/node_Exporter/templates/nodeExporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User={{ user }}
ExecStart={{ node_exporter_home }}/node_exporter

[Install]
WantedBy=default.target
